# supervisord-conf/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:postfix]
process_name	= postfix
directory	= /etc/postfix
command		= /usr/sbin/postfix -c /etc/postfix start-fg -v
startsecs	= 0
autorestart	= true